# runit/config/etc/sv/agetty-tty6/run
#!/bin/sh
exec setsid /sbin/agetty 38400 tty6 linux
