# runit/config/etc/sv/agetty-tty3/run
#!/bin/sh
exec setsid /sbin/agetty 38400 tty3 linux
